# pdfmaster-updater.service
[Unit]
Description=PDF Master Auto Updater
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-update.sh
StandardOutput=journal
StandardError=journal

# 安全設定
User=root
PrivateTmp=true
NoNewPrivileges=true

# 環境變數（固定配置）
Environment="TELEGRAM_BOT_TOKEN=7002177842:AAE7cXJpmqXKmZAh19aef2P4dfnxW0GVjZs"
Environment="TELEGRAM_CHAT_ID=7341258916" 